# api/supervisor.conf
[program:timelog]
command = /usr/bin/python timelog-api.py
directory = /service
stdout_logfile = /timelog-api.log
