# api/supervisor.conf
[program:timelog]
command = python app.py
directory = /service
environment=USER=root,PYTHONPATH=/usr/bin/
stdout_logfile = /timelog-out.log
stderr_logfile = /timelog-err.log
